# ryzen/autostart/picom.desktop
[Desktop Entry]
Encoding=UTF-8
Type=Application
Name=Picom
Comment=X11 compositor
Exec=picom -b
StartupNotify=false
Terminal=false
Hidden=false
